# cherryscrapper.conf
[program:cherryscrapper]
command = /home/ubuntu/www/cherryscrapper/scrapper.sh
stdout_logfile=/home/ubuntu/www/supervisor-cherryscrapper.log
autostart=true
autorestart=true